# huming17/DZMVC-1.0-DWZ-Example/source/module/crontab/cron_notify_trigger.sh
#!/bin/sh
wget http://172.16.1.7:9002/index.php?mod=crontab&action=site_effective --output-document=/tmp/cron_temp.html